# roles/kibana/templates/etc/systemd/system/kibana.service.j2
[Unit]
Description=Kibana 3

[Service]
Type=simple
Environment=CONFIG_PATH=/opt/kibana/config/kibana.yml
Environment=NODE_ENV=production
ExecStart=/opt/kibana/bin/kibana

[Install]
WantedBy=multi-user.target